# battery_notify/battery-notify.service
[Unit]
Description=Battery charge notifier

[Service]
ExecStart=/usr/local/bin/battery-notify.sh

[Install]
WantedBy=default.target
